# ai-stack/vtuber-tvs99-manager/vtuber-tvs99-manager.service
[Unit]
Description=VTuber TVS99 Manager Stack
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/media/rizzo/RAIDSTATION/stacks/ai-stack/vtuber-tvs99-manager
ExecStart=/media/rizzo/RAIDSTATION/stacks/ai-stack/vtuber-tvs99-manager/manage-stack.sh start
ExecStop=/media/rizzo/RAIDSTATION/stacks/ai-stack/vtuber-tvs99-manager/manage-stack.sh stop
ExecReload=/media/rizzo/RAIDSTATION/stacks/ai-stack/vtuber-tvs99-manager/manage-stack.sh restart
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
